# app.service
[Unit]
Description=Flask Application Service
After=network.target

[Service]
User=ubuntu
ExecStart=/usr/local/bin/gunicorn -b 0.0.0.0:8000 app:app
WorkingDirectory=/path/to/app
Restart=always

[Install]
WantedBy=multi-user.target
